# paquete_listo_latest_10_08_2025/DEBIAN/postinst
#!/bin/sh
# postinst: se ejecuta después de la instalación del paquete.

# --- SECCIÓN DE PERMISOS AGREGADA ---

# Define el usuario que ejecutará el servicio.
# ¡IMPORTANTE! Este usuario debe ser el mismo que configures en tu archivo .service
# Un usuario común para servicios web es 'www-data'.
APP_USER="root" # <-- ¡CAMBIA ESTO SI ES NECESARIO!
APP_DIR="/opt/mikrotik-manager"

echo "Estableciendo propietario de los archivos a '$APP_USER'..."
# La opción -R (recursivo) aplica los permisos a la carpeta y todo su contenido.
chown -R $APP_USER:$APP_USER $APP_DIR
chmod -R u+rwX $APP_DIR

# --- FIN DE LA SECCIÓN DE PERMISOS ---


echo "Recargando systemd..."
systemctl daemon-reload

echo "Habilitando el servicio mikrotik-manager..."
systemctl enable mikrotik-manager.service

echo "Iniciando el servicio mikrotik-manager..."
systemctl start mikrotik-manager.service

echo "La interfaz web debería estar disponible en http://<ip-del-servidor>:8080"

exit 0